# etc/wsgi/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=info                                   ; info, debug, warn, trace, error
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; run supervisord NOT as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisord/               ; where child log files will live

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

; This is where we run individual Aiohttp app instances.

[program:aiohttp-8000]
command=/usr/local/bin/python /app/run.py --port=8000
stderr_logfile = /var/log/supervisord/aiohttp_err.log
stdout_logfile = /var/log/supervisord/aiohttp_out.log

[program:aiohttp-8001]
command=/usr/local/bin/python /app/run.py --port=8001
stderr_logfile = /var/log/supervisord/aiohttp_err.log
stdout_logfile = /var/log/supervisord/aiohttp_out.log

[program:aiohttp-8002]
command=/usr/local/bin/python /app/run.py --port=8002
stderr_logfile = /var/log/supervisord/aiohttp_err.log
stdout_logfile = /var/log/supervisord/aiohttp_out.log

[program:aiohttp-8003]
command=/usr/local/bin/python /app/run.py --port=8003
stderr_logfile = /var/log/supervisord/aiohttp_err.log
stdout_logfile = /var/log/supervisord/aiohttp_out.log
